# dist/supervisord.ini
[supervisord]
nodaemon=true

[program:rsyslog]
autorestart=true
command=rsyslogd -n
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root

[program:cron]
autorestart=true
command=crond -f
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root

[program:nginx]
autorestart=true
command=nginx -g "daemon off;"
stdout_logfile=syslog
stderr_logfile=syslog
username=%(ENV_USR_USER)s
