# dist/supervisord.ini
[supervisord]
nodaemon=true

[program:nginx]
autorestart=true
command=nginx
directory=%(MY_WORKDIR_ENV_VAR)s
stdout_logfile=syslog
stderr_logfile=syslog
user=%(MY_USER_ENV_VAR)s

[program:cron]
autorestart=true
command=crond -f
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root

[program:rsyslog]
autorestart=true
command=rsyslogd -n
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root
